# package/disp-settings/src/disp-settings-dual-display-restore.service.in
[Unit]
Description=Restore disp-settings dual-display brightness
Wants=als-dimmer.service
After=local-fs.target als-dimmer.service
Before=qt-demo-launcher.service

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/disp-settings-dual-display-restore.py
TimeoutStartSec=25
StandardOutput=journal
StandardError=journal
SyslogIdentifier=disp-settings-dual-display-restore

[Install]
WantedBy=multi-user.target
