# board/pi4-touch-demo/fs-overlay/usr/lib/systemd/system/systemd-himaxtouch.service
[Unit]
Description=himax-touch-driver

[Service]
Type=oneshot
ExecStart=/etc/init.d/S59TouchBlinkServer start
ExecStop=/etc/init.d/S59TouchBlinkServer stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
